# sd_fs/binary_developer/var/lib/dpkg/info/lm-sensors.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lm-sensors.service
/usr
/usr/bin
/usr/bin/sensors
/usr/bin/sensors-conf-convert
/usr/sbin
/usr/sbin/sensors-detect
/usr/share
/usr/share/doc
/usr/share/doc/lm-sensors
/usr/share/doc/lm-sensors/applications
/usr/share/doc/lm-sensors/temperature-sensors
/usr/share/doc/lm-sensors/fan-divisors.gz
/usr/share/doc/lm-sensors/copyright
/usr/share/doc/lm-sensors/README.Debian.gz
/usr/share/doc/lm-sensors/examples
/usr/share/doc/lm-sensors/examples/daemon
/usr/share/doc/lm-sensors/examples/daemon/healthd.sh
/usr/share/doc/lm-sensors/examples/hotplug
/usr/share/doc/lm-sensors/examples/hotplug/unhide_ICH_SMBus
/usr/share/doc/lm-sensors/examples/hotplug/README.p4b
/usr/share/doc/lm-sensors/examples/tellerstats
/usr/share/doc/lm-sensors/examples/tellerstats/README
/usr/share/doc/lm-sensors/examples/tellerstats/gather.sh
/usr/share/doc/lm-sensors/examples/tellerstats/tellerstats.conf
/usr/share/doc/lm-sensors/examples/tellerstats/gnuplotscript.tmpl
/usr/share/doc/lm-sensors/examples/tellerstats/index.shtml
/usr/share/doc/lm-sensors/examples/tellerstats/tellerstats.sh
/usr/share/doc/lm-sensors/vid
/usr/share/doc/lm-sensors/donations.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/sensors-conf-convert.8.gz
/usr/share/man/man8/sensors-detect.8.gz
/usr/share/man/man5
/usr/share/man/man5/sensors.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/sensors.1.gz
/etc
/etc/init.d
/etc/init.d/lm-sensors
/usr/share/doc/lm-sensors/changelog.Debian.gz
