# sd_fs/binary_developer/var/lib/dpkg/info/network-manager.list
/.
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/org.freedesktop.NetworkManager.pkla
/var/lib/NetworkManager
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/nm-dispatcher.conf
/etc/dbus-1/system.d/org.freedesktop.NetworkManager.conf
/etc/dbus-1/system.d/nm-ofono.conf
/etc/dbus-1/system.d/nm-avahi-autoipd.conf
/etc/init
/etc/init/network-manager.conf
/etc/dnsmasq.d
/etc/dnsmasq.d/network-manager
/etc/init.d
/etc/init.d/network-manager
/etc/NetworkManager
/etc/NetworkManager/dnsmasq.d
/etc/NetworkManager/system-connections
/etc/NetworkManager/conf.d
/etc/NetworkManager/NetworkManager.conf
/etc/NetworkManager/VPN
/etc/NetworkManager/dispatcher.d
/etc/NetworkManager/dispatcher.d/pre-up.d
/etc/NetworkManager/dispatcher.d/01ifupdown
/etc/NetworkManager/dispatcher.d/pre-down.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/NetworkManager-wait-online.service
/lib/systemd/system/NetworkManager.service
/lib/systemd/system/NetworkManager-dispatcher.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/77-nm-olpc-mesh.rules
/usr
/usr/lib
/usr/lib/pppd
/usr/lib/pppd/2.4.6
/usr/lib/pppd/2.4.6/nm-pppd-plugin.so
/usr/lib/NetworkManager
/usr/lib/NetworkManager/nm-avahi-autoipd.action
/usr/lib/NetworkManager/nm-dhcp-helper
/usr/lib/NetworkManager/ifblacklist_migrate.sh
/usr/lib/NetworkManager/nm-dispatcher
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/NetworkManager
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-wwan.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-settings-plugin-ofono.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-wwan.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-settings-plugin-ifupdown.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-bluetooth.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-adsl.so
/usr/lib/arm-linux-gnueabihf/NetworkManager/libnm-device-plugin-wifi.so
/usr/bin
/usr/bin/nmcli
/usr/bin/nm-online
/usr/bin/nmtui
/usr/sbin
/usr/sbin/NetworkManager
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/NetworkManager.8.gz
/usr/share/man/man5
/usr/share/man/man5/NetworkManager.conf.5.gz
/usr/share/man/man5/nmcli-examples.5.gz
/usr/share/man/man5/nm-system-settings.conf.5.gz
/usr/share/man/man5/nm-settings.5.gz
/usr/share/man/man1
/usr/share/man/man1/nmcli.1.gz
/usr/share/man/man1/nm-online.1.gz
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/nmcli
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_network-manager.py
/usr/share/polkit-1
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/60-network-manager.rules
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
/usr/share/doc
/usr/share/doc/network-manager
/usr/share/doc/network-manager/examples
/usr/share/doc/network-manager/examples/server.conf
/usr/share/doc/network-manager/copyright
/usr/share/doc/network-manager/README.gz
/usr/share/doc/network-manager/README.Debian
/usr/share/doc/network-manager/NEWS.Debian.gz
/usr/share/doc/network-manager/AUTHORS
/usr/share/doc/network-manager/NEWS.gz
/lib/systemd/system/network-manager.service
/usr/bin/nmtui-hostname
/usr/bin/nmtui-edit
/usr/bin/nmtui-connect
/usr/share/apport/package-hooks/source_network-manager-applet.py
/usr/share/doc/network-manager/changelog.Debian.gz
