# sd_fs/binary_developer/var/lib/dpkg/info/xen-utils-common.list
/.
/var
/var/lib
/var/lib/xen
/etc
/etc/xen
/etc/xen/xend-pci-quirks.sxp
/etc/xen/xend-pci-permissive.sxp
/etc/xen/scripts
/etc/xen/scripts/vif-openvswitch
/etc/xen/scripts/vif2
/etc/xen/scripts/hotplugpath.sh
/etc/xen/scripts/vscsi
/etc/xen/scripts/vif-common.sh
/etc/xen/scripts/block-nbd
/etc/xen/scripts/block-enbd
/etc/xen/scripts/block-common.sh
/etc/xen/scripts/block
/etc/xen/scripts/xen-network-common.sh
/etc/xen/scripts/vif-route
/etc/xen/scripts/block-drbd-probe
/etc/xen/scripts/locking.sh
/etc/xen/scripts/external-device-migrate
/etc/xen/scripts/logging.sh
/etc/xen/scripts/block-iscsi
/etc/xen/scripts/xen-hotplug-cleanup
/etc/xen/scripts/vif-bridge
/etc/xen/scripts/vif-nat
/etc/xen/scripts/qemu-ifup
/etc/xen/scripts/xen-script-common.sh
/etc/xen/scripts/xen-hotplug-common.sh
/etc/xen/scripts/vif-setup
/etc/xen/xlexample.hvm
/etc/xen/xl.conf
/etc/xen/xend-config.sxp
/etc/xen/xlexample.pvlinux
/etc/default
/etc/default/xendomains
/etc/default/xend
/etc/init.d
/etc/init.d/xen
/etc/init.d/xendomains
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/xen-backend.rules
/usr
/usr/share
/usr/share/doc
/usr/share/doc/xen-utils-common
/usr/share/doc/xen-utils-common/changelog.Debian.gz
/usr/share/doc/xen-utils-common/README.Debian
/usr/share/doc/xen-utils-common/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/xentrace.8.gz
/usr/share/man/man5
/usr/share/man/man5/xlcpupool.cfg.5.gz
/usr/share/man/man5/xl.cfg.5.gz
/usr/share/man/man5/xl.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/xenstore-chmod.1.gz
/usr/share/man/man1/xenstore.1.gz
/usr/share/man/man1/xentrace_format.1.gz
/usr/share/man/man1/xentop.1.gz
/usr/share/man/man1/xl.1.gz
/usr/share/man/man1/xenstore-ls.1.gz
/usr/share/xen-utils-common
/usr/share/xen-utils-common/default.xen
/usr/sbin
/usr/sbin/xen
/usr/lib
/usr/lib/xen-common
/usr/lib/xen-common/bin
/usr/lib/xen-common/bin/xen-toolstack-wrapper
/usr/lib/xen-common/bin/xen-dir
/usr/lib/xen-common/bin/xen-utils-wrapper
/usr/lib/xen-common/bin/xen-toolstack
/usr/lib/xen-common/bin/xen-init-name
/usr/lib/xen-common/bin/xen-init-list
/usr/lib/xen-common/bin/xen-version
/usr/sbin/xentrace_format
/usr/sbin/xentop
/usr/sbin/xentrace
/usr/sbin/xl
/usr/sbin/xm
/usr/sbin/xenperf
/usr/sbin/xentrace_setmask
/usr/sbin/xentrace_setsize
/usr/sbin/xenpm
